# install.soulmate-transmute/usr/share/soulmate/profile.d/soulmate-default-vars.sh
TARGET_RELEASE="9"
SM_DEFMIRROR="http://cloud.ausiasmarch.net/soulmate/"
DEBIAN_CODENAME="stretch"
DEBIAN_DEFMIRROR="http://ftp.debian.org/debian/"
DEBIAN_DEFSECURITY="http://security.debian.org"

