# install.transmute/usr/share/aurex/profile.d/aurex-default-vars.sh
AUREX_RELEASE="12.04"
UBUNTU_CODENAME="precise"
AUREX_CODENAME="precise"
UBUNTU_DEFMIRROR="http://archive.ubuntu.com/ubuntu/"
AUREX_DEFMIRROR="http://intranet.ausiasmarch.net/aurex/"
